# conf/chilli.in
#!/bin/sh
#
# chilli      CoovaChilli init
#
# chkconfig: 2345 65 35
# description: CoovaChilli

# Source function library.

daemon=
killproc=killall

[ -e /etc/rc.d/init.d/functions ] && {
    . /etc/rc.d/init.d/functions
    daemon=daemon
    killproc=killproc
}

[ -e /etc/sysconfig/network ] && {
    . /etc/sysconfig/network
    [ ${NETWORKING} = "no" ] && exit 0
}

[ -f @SBINDIR@/chilli ] || exit 0

. @ETCCHILLI@/functions

MULTI=$(ls @ETCCHILLI@/*/chilli.conf 2>/dev/null)
[ -z "$DHCPIF" ] && [ -n "$MULTI" ] && {
    for c in $MULTI; 
    do
	echo "Found configuration $c"
	DHCPIF=$(basename $(echo $c|sed 's#/chilli.conf##'))
	export DHCPIF
	echo "Running DHCPIF=$DHCPIF $0 $*"
	sh $0 $*
    done
    exit
}

if [ -n "$DHCPIF" ]; then
    CONFIG=@ETCCHILLI@/$DHCPIF/chilli.conf
else
    CONFIG=@SYSCONFDIR@/chilli.conf
fi

[ -f $CONFIG ] || {
    echo "$CONFIG Not found"
    exit 0
}

check_required

RETVAL=0
prog="chilli"

case $1 in
    start)
        echo -n "Starting $prog: "

        /sbin/modprobe tun >/dev/null 2>&1
        echo 1 > /proc/sys/net/ipv4/ip_forward

	[ -e /dev/net/tun ] || {
	    (cd /dev; 
		mkdir net; 
		cd net; 
		mknod tun c 10 200)
	}

	writeconfig
        radiusconfig

	test ${HS_ADMINTERVAL:-0} -gt 0 && {	
            (crontab -l 2>&- | grep -v $0
		echo "*/$HS_ADMINTERVAL * * * * $0 radconfig"
		) | crontab - 2>&-
	}
	
	ifconfig $HS_LANIF 0.0.0.0
        $daemon @SBINDIR@/chilli -c $CONFIG &
        echo
	;;
    
    radconfig)
	[ -e $MAIN_CONF ] || writeconfig
	radiusconfig
	;;

    reload)
	killall -HUP chilli
	;;

    restart)
	$0 stop
	sleep 1
	$0 start
	RETVAL=$?
	;;
    
    stop)
        echo -n $"Shutting down $prog: "

	crontab -l 2>&- | grep -v $0 | crontab -

        [ -f @VARRUN@/chilli.$HS_LANIF.pid ] && { 
	    kill $(cat @VARRUN@/chilli.$HS_LANIF.pid)
	    RETVAL=$
	    [ $RETVAL = 0 ] && {
		rm -f @VARRUN@/chilli.$HS_LANIF.pid 2>/dev/null
	    }
	}

	#$killproc chilli

        echo
	;;
    
    *)
        echo "Usage: $0 {start|stop|restart|reload|radconfig}"
        exit 1
esac

exit $?
